# lamp_voice.conf
[program:lamp_voice]
directory=/home/pi/lampi-voice
command=python3 /home/pi/lampi-voice/voiceLampi.py
environment=HOME="/home/pi"
priority=500
user=pi
autostart=true
autorestart=true